# etc/supervisor/conf.d/rpi_bridge_bureau.conf
[program:redis_serial_sync_internet]
command=redis_serial_sync -p /dev/ttyUSB0
user=pi
directory=/home/pi/
autostart=true
autorestart=true
environment=PYTHONUNBUFFERED="1"
